# deploy/com.stackowl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stackowl.stackowl</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/stackowl</string>
        <string>serve</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/StackOwl/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/StackOwl/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>STACKOWL_HOME</key>
        <string>~/Library/Application Support/StackOwl</string>
    </dict>
</dict>
</plist>
